# udev/rules.d/70-snap.sonobus.rules
# This file is automatically generated.
# alsa
KERNEL=="controlC[0-9]*", TAG+="snap_sonobus_sonobus"
# alsa
KERNEL=="hwC[0-9]*D[0-9]*", TAG+="snap_sonobus_sonobus"
# alsa
KERNEL=="midiC[0-9]*D[0-9]*", TAG+="snap_sonobus_sonobus"
# alsa
KERNEL=="pcmC[0-9]*D[0-9]*[cp]", TAG+="snap_sonobus_sonobus"
# alsa
KERNEL=="seq", TAG+="snap_sonobus_sonobus"
# alsa
KERNEL=="timer", TAG+="snap_sonobus_sonobus"
# alsa
SUBSYSTEM=="sound", KERNEL=="card[0-9]*", TAG+="snap_sonobus_sonobus"
TAG=="snap_sonobus_sonobus", SUBSYSTEM!="module", SUBSYSTEM!="subsystem", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_sonobus_sonobus $devpath $major:$minor"
