# systemd/system/var-lib-snapd-snap-snap\x2dstore-1270.mount
[Unit]
Description=Mount unit for snap-store, revision 1270
After=snapd.mounts-pre.target
Before=snapd.mounts.target

[Mount]
What=/var/lib/snapd/snaps/snap-store_1270.snap
Where=/var/lib/snapd/snap/snap-store/1270
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
